# puppet/modules/lib_batchs/templates/noids.conf.erb
description "Noids Unique identifier service"

start on runlevel [2345]
console none
respawn

exec <%= @noid_root %>/bin/noids \
     --config  <%= @noid_root %>/config.ini \
     --log <%= @noid_root %>/log/noids.log \
     --pid <%= @noid_root %>/noids.pid
